# Jeopardy-style/challenge/Reverse/rpi-door/fs/var/lib/dpkg/info/ifplugd.list
/.
/usr
/usr/sbin
/usr/sbin/ifplugstatus
/usr/sbin/ifplugd
/usr/sbin/ifstatus
/usr/share
/usr/share/doc
/usr/share/doc/ifplugd
/usr/share/doc/ifplugd/NEWS.gz
/usr/share/doc/ifplugd/FAQ
/usr/share/doc/ifplugd/README.gz
/usr/share/doc/ifplugd/SUPPORTED_DRIVERS
/usr/share/doc/ifplugd/changelog.Debian.gz
/usr/share/doc/ifplugd/style.css
/usr/share/doc/ifplugd/README.html
/usr/share/doc/ifplugd/8139too.c.0.9.26.patch
/usr/share/doc/ifplugd/changelog.gz
/usr/share/doc/ifplugd/copyright
/usr/share/bug
/usr/share/bug/ifplugd
/usr/share/bug/ifplugd/script
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ifplugd.8.gz
/usr/share/man/man8/ifplugstatus.8.gz
/usr/share/man/man5
/usr/share/man/man5/ifplugd.conf.5.gz
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-ifplugd.rules
/lib/udev/ifplugd.agent
/etc
/etc/pm
/etc/pm/sleep.d
/etc/pm/sleep.d/74ifplugd
/etc/init.d
/etc/init.d/ifplugd
/etc/apm
/etc/apm/scripts.d
/etc/apm/scripts.d/ifplugd
/etc/apm/other.d
/etc/apm/resume.d
/etc/apm/suspend.d
/etc/apm/event.d
/etc/ifplugd
/etc/ifplugd/ifplugd.action
/etc/ifplugd/action.d
/etc/ifplugd/action.d/ifupdown
/usr/share/man/man8/ifstatus.8.gz
